# rootfs/etc/services.d/nzbhydra2/run
#!/usr/bin/execlineb -P

with-contenv
multisubstitute
{
    importas PUID PUID
    importas PGID PGID
}
s6-setuidgid ${PUID}:${PGID}
./nzbhydra2 --datafolder /config
